# updatedns/updatedns.service
[Unit]
Description=UpdateDNS

[Service]
Type=simple
ExecStart=updatedns --monitor

[Install]
WantedBy=multi-user.target

