# updatedns.service
[Unit]
Description=Update DNS

[Service]
Type=simple
ExecStart=/opt/updatedns/updatedns

[Install]
WantedBy=multi-user.target

